# docker/controller/supervisord_controller_files/cassandra.ini
[program:cassandra]
command=cassandra -f
stdout_logfile = /var/log/cassandra/cassandra.out
stderr_logfile = /var/log/cassandra/cassandra.err
autorestart = true
stopasgroup=true
