# docker/controller/supervisord_controller_files/rabbitmq.ini
[program:rabbitmq]
command=/rabbitmq-launch.sh
stdout_logfile = /var/log/rabbitmq/rabbitmq.out
stderr_logfile = /var/log/rabbitmq/rabbitmq.err
autorestart = true
stopasgroup=true
