# docker/controller/supervisord_controller_files/redis.ini
[program:redis]
command=/usr/bin/redis-server
stdout_logfile = /var/log/redis/redis-server.out
stderr_logfile = /var/log/redis/redis-server.err
autorestart = true
stopasgroup=true
